# crontab
@reboot sleep 30 && cd /home/pi/plants/ && sudo git pull && ./webpage-updater.sh
0 22 * * * /home/pi/plants/timelapse.sh
